# luci-app-aurora-config/root/etc/uci-defaults/80_aurora
#!/bin/sh

readonly CONF_FILE="/etc/config/aurora"
readonly CLASSIC_TEMPLATE="/usr/share/aurora/classic.template"
readonly SAGE_GREEN_TEMPLATE="/usr/share/aurora/sage-green.template"
readonly AMBER_SAND_TEMPLATE="/usr/share/aurora/amber-sand.template"
readonly MONOCHROME_TEMPLATE="/usr/share/aurora/monochrome.template"
readonly SKY_BLUE_TEMPLATE="/usr/share/aurora/sky-blue.template"
readonly LOG_TAG="aurora"

# Bump when the meaning of stored color tokens changes incompatibly.
# v2: two-layer flip -- light_/dark_ status inputs (success/info/warning/danger)
#     went from "light surface" to "dark content", surfaces derived via
#     oklch(from var(--token) ...). Pre-flip snapshots leave content the same
#     colour as the derived surface -> invisible label text and icons.
# v3: fixes a re-seed bug in v2 that left ~14 colour keys empty (a nested
#     command substitution inside a pipe-fed `while read` swallowed stdin and
#     skipped keys). v3 re-runs a corrected re-seed, so devices stuck at v2
#     with a half-migrated config are repaired.
# v4: single-plane token redesign. The two-layer taxonomy (canvas/content/
#     surface_raised + ~57 derived expression tokens) is gone. The schema is now
#     10 editable inputs (canvas->bg, content->text, surface_raised dropped) plus
#     17 engine-derived flat colours. The re-seed prunes EVERY legacy light_/dark_
#     key absent from the new template, clearing stale overrides that header.ut
#     would otherwise inject against now-nonexistent CSS variables. Custom colours
#     reset to the active template defaults (no value mapping).
readonly SCHEMA_VERSION="4"

log() {
    logger -t "$LOG_TAG" "$1"
}

# Re-seed every light_/dark_ colour token from the template. Key lists are
# captured into variables first so no uci command runs inside a pipe-fed
# `while read` (a nested command substitution there can swallow stdin and
# silently skip keys). Sets first, then prunes legacy keys, so no key is ever
# left deleted-but-unset. Fonts (struct_*), toolbar items, nav mode, icon and
# logo/favicon options are left untouched.
reseed_colors_from_template() {
    local template="$1" dir key value tmpl_keys live_keys

    dir=$(mktemp -d /tmp/aurora-migrate.XXXXXX) || return 1
    cp "$template" "$dir/aurora" || { rm -rf "$dir"; return 1; }

    tmpl_keys=$(uci -c "$dir" -q show aurora.theme \
        | sed -n 's/^aurora\.theme\.\(\(light\|dark\)_[^=]*\)=.*/\1/p')
    live_keys=$(uci -q show aurora.theme \
        | sed -n 's/^aurora\.theme\.\(\(light\|dark\)_[^=]*\)=.*/\1/p')

    [ -n "$tmpl_keys" ] || { rm -rf "$dir"; return 1; }

    for key in $tmpl_keys; do
        value=$(uci -c "$dir" -q get "aurora.theme.$key") || { rm -rf "$dir"; return 1; }
        uci -q set "aurora.theme.$key=$value" || { rm -rf "$dir"; return 1; }
    done

    for key in $live_keys; do
        uci -c "$dir" -q get "aurora.theme.$key" >/dev/null 2>&1 \
            || uci -q delete "aurora.theme.$key"
    done

    rm -rf "$dir"
    return 0
}

migrate_color_schema() {
    local template="$1" stored

    stored=$(uci -q get aurora.theme.config_version)
    [ "$stored" = "$SCHEMA_VERSION" ] && return 0

    log "Migrating colour tokens to schema v$SCHEMA_VERSION (was '${stored:-unset}')"
    if reseed_colors_from_template "$template"; then
        uci -q set aurora.theme.config_version="$SCHEMA_VERSION"
        uci -q commit aurora
        log "Colour token migration complete"
    else
        log "Colour token migration failed; leaving existing config untouched"
    fi
}

main() {
    log "Initializing Aurora theme configuration"

    template=""
    if [ -f "$CLASSIC_TEMPLATE" ]; then
        template="$CLASSIC_TEMPLATE"
    elif [ -f "$SAGE_GREEN_TEMPLATE" ]; then
        template="$SAGE_GREEN_TEMPLATE"
    elif [ -f "$AMBER_SAND_TEMPLATE" ]; then
        template="$AMBER_SAND_TEMPLATE"
    elif [ -f "$MONOCHROME_TEMPLATE" ]; then
        template="$MONOCHROME_TEMPLATE"
    elif [ -f "$SKY_BLUE_TEMPLATE" ]; then
        template="$SKY_BLUE_TEMPLATE"
    fi

    if [ -z "$template" ]; then
        log "Template missing: $CLASSIC_TEMPLATE, $SAGE_GREEN_TEMPLATE, $AMBER_SAND_TEMPLATE, $MONOCHROME_TEMPLATE, $SKY_BLUE_TEMPLATE"
        exit 0
    fi

    if [ ! -f "$CONF_FILE" ]; then
        if cp "$template" "$CONF_FILE"; then
            chmod 644 "$CONF_FILE"
            uci -q set aurora.theme.config_version="$SCHEMA_VERSION"
            uci -q commit aurora
            log "First-time initialization: Config created from template at schema v$SCHEMA_VERSION"
        else
            log "Failed to create config from template"
        fi
    else
        log "Configuration already exists, checking colour token schema"
        migrate_color_schema "$template"
    fi
}

main

exit 0
